# common/config/supervisord/conf.d/dbus.conf
[program:dbus-daemon]
command=/usr/bin/dbus-daemon --system --nofork
